# sysutils/freebsd_controller/files/freebsd_controller.in
#!/bin/sh

# PROVIDE: freebsd_controller
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# freebsd_controller_enable (bool):          Set to NO by default.
#               Set it to YES to enable freebsd_controller.
# freebsd_controller_user (string):          Set user that freebsd_controller will run under
#               Default is "root".
# freebsd_controller_group (string):         Set group that freebsd_controller will run under
#               Default is "wheel".
# freebsd_controller_args (string):          Set extra arguments to pass to freebsd_controller
#               Default is "".
# freebsd_controller_listen_address (string):Set ip:port that freebsd_controller will listen on
#               Default is ":9100".

. /etc/rc.subr

name=freebsd_controller
rcvar=freebsd_controller_enable

load_rc_config $name

: ${freebsd_controller_enable:="NO"}
: ${freebsd_controller_user:="root"}
: ${freebsd_controller_group:="wheel"}
: ${freebsd_controller_args:=""}


pidfile=/var/run/freebsd_controller.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/freebsd_controller"
command_args="-S -T ${name} -f -p ${pidfile} /usr/bin/env PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin ${procname} \
    ${freebsd_controller_args}"

start_precmd=freebsd_controller_startprecmd

freebsd_controller_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${freebsd_controller_user} -g ${freebsd_controller_group} /dev/null ${pidfile};
    fi
}

load_rc_config $name
run_rc_command "$1"
