# net/arp_exporter/files/arp_exporter.in
#!/bin/sh

# PROVIDE: arp_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# arp_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable arp_exporter.
# arp_exporter_user (string):          Set user that arp_exporter will run under
#               Default is "nobody".
# arp_exporter_group (string):         Set group that arp_exporter will run under
#               Default is "nobody".
# arp_exporter_args (string):          Set extra arguments to pass to arp_exporter
#               Default is "".
# arp_exporter_listen_address (string):Set ip:port that arp_exporter will listen on
#               Default is ":9100".

. /etc/rc.subr

name=arp_exporter
rcvar=arp_exporter_enable

load_rc_config $name

: ${arp_exporter_enable:="NO"}
: ${arp_exporter_user:="nobody"}
: ${arp_exporter_group:="nobody"}
: ${arp_exporter_args:=""}
: ${arp_exporter_textfile_dir:="/var/tmp/arp_exporter"}


pidfile=/var/run/arp_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/arp_exporter"
command_args="-f -p ${pidfile} /usr/bin/env ${procname} \
    ${arp_exporter_args}"

start_precmd=arp_exporter_startprecmd

arp_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${arp_exporter_user} -g ${arp_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${arp_exporter_textfile_dir} ]; then
        install -d -o ${arp_exporter_user} -g ${arp_exporter_group} -m 1755 ${arp_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
