# sysutils/pf_exporter/files/pf_exporter.in
#!/bin/sh

# PROVIDE: pf_exporter
# REQUIRE: LOGIN pf
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# pf_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable pf_exporter.
# pf_exporter_args (string):          Set extra arguments to pass to pf_exporter
#               Default is "".
# pf_exporter_listen_address (string):Set ip:port that pf_exporter will listen on
#               Default is ":9107".

. /etc/rc.subr

name=pf_exporter
rcvar=pf_exporter_enable

load_rc_config $name

: ${pf_exporter_enable:="NO"}
: ${pf_exporter_args:=""}
: ${pf_exporter_listen_address:=":9107"}

_pidfile=/var/run/${name}.pid
_binary="%%PREFIX%%/bin/pf_exporter"

start_cmd="${name}_start"
stop_cmd="${name}_stop"
status_cmd="${name}_status"

pf_exporter_start()
{
	/usr/sbin/daemon -S -T ${name} -f -p ${_pidfile} \
	    ${_binary} -web.listen-address ${pf_exporter_listen_address} \
	    ${pf_exporter_args}
}

pf_exporter_stop()
{
	if [ -f ${_pidfile} ]; then
		kill $(cat ${_pidfile}) && rm -f ${_pidfile}
	else
		echo "${name} not running."
	fi
}

pf_exporter_status()
{
	if [ -f ${_pidfile} ] && kill -0 $(cat ${_pidfile}) 2>/dev/null; then
		echo "${name} is running as pid $(cat ${_pidfile})."
	else
		echo "${name} is not running."
		return 1
	fi
}

run_rc_command "$1"
