# comms/znet/files/agent.in
#!/bin/sh

# PROVIDE: agent
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# agent_enable (bool):          Set to NO by default.
#               Set it to YES to enable agent.
# agent_user (string):          Set user that agent will run under
#               Default is "nobody".
# agent_group (string):         Set group that agent will run under
#               Default is "nobody".
# agent_args (string):          Set extra arguments to pass to agent
#               Default is "".

. /etc/rc.subr

name="agent"
rcvar="agent_enable"
pidfile=/var/run/agent.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"
command_args="-S -T ${name} -r -f -P ${pidfile} /usr/bin/env ${procname} ${agent_args}"

load_rc_config $name

: ${agent_enable:="NO"}
: ${agent_user:="nobody"}
: ${agent_group:="nobody"}
: ${agent_args:="agent --config /usr/local/etc/agent.yaml"}

run_rc_command "$1"
