# comms/znet/files/agent.in
#!/bin/sh

# PROVIDE: agent
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# agent_enable (bool):          Set to NO by default.
#               Set it to YES to enable agent.
# agent_user (string):          Set user that agent will run under
#               Default is "nobody".
# agent_group (string):         Set group that agent will run under
#               Default is "nobody".
# agent_args (string):          Set extra arguments to pass to agent
#               Default is "".

. /etc/rc.subr

name="agent"
rcvar="agent_enable"
pidfile="/var/run/${name}.pid"
supervisor_pidfile="/var/run/supervisor_${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${agent_enable:="NO"}
: ${agent_user:="nobody"}
: ${agent_group:="nobody"}
: ${agent_args:="${name} --config /usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} -P ${supervisor_pidfile} ${procname} ${agent_args}"

start_precmd=agent_startprecmd

agent_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${agent_user} -g ${agent_group} /dev/null ${pidfile};
    fi

    if [ ! -e ${supervisor_pidfile} ]; then
        install -o ${agent_user} -g ${agent_group} /dev/null ${supervisor_pidfile};
    fi
}

run_rc_command "$1"
