# comms/znet/files/arpwatch.in
#!/bin/sh

# PROVIDE: arpwatch
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# arpwatch_enable (bool):          Set to NO by default.
#               Set it to YES to enable arpwatch.
# arpwatch_user (string):          Set user that arpwatch will run under
#               Default is "nobody".
# arpwatch_group (string):         Set group that arpwatch will run under
#               Default is "nobody".
# arpwatch_args (string):          Set extra arguments to pass to arpwatch
#               Default is "".

. /etc/rc.subr

name="arpwatch"
rcvar="arpwatch_enable"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${arpwatch_enable:="NO"}
: ${arpwatch_user:="nobody"}
: ${arpwatch_group:="nobody"}
: ${arpwatch_args:="${name} --config /usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} ${procname} ${arpwatch_args}"

run_rc_command "$1"
