# comms/znet/files/builder.in
#!/bin/sh

# PROVIDE: builder
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# builder_enable (bool):          Set to NO by default.
#               Set it to YES to enable builder.
# builder_user (string):          Set user that builder will run under
#               Default is "nobody".
# builder_group (string):         Set group that builder will run under
#               Default is "nobody".
# builder_args (string):          Set extra arguments to pass to builder
#               Default is "".

. /etc/rc.subr

name="builder"
rcvar="builder_enable"
pidfile="/var/run/${name}.pid"
supervisor_pidfile="/var/run/supervisor_${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${builder_enable:="NO"}
: ${builder_user:="nobody"}
: ${builder_group:="nobody"}
: ${builder_args:="${name} --config /usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} -P ${supervisor_pidfile} ${procname} ${builder_args}"

run_rc_command "$1"
