# comms/znet/files/gitwatch.in
#!/bin/sh

# PROVIDE: gitwatch
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# gitwatch_enable (bool):          Set to NO by default.
#               Set it to YES to enable gitwatch.
# gitwatch_user (string):          Set user that gitwatch will run under
#               Default is "nobody".
# gitwatch_group (string):         Set group that gitwatch will run under
#               Default is "nobody".
# gitwatch_args (string):          Set extra arguments to pass to gitwatch
#               Default is "".

. /etc/rc.subr

name="gitwatch"
rcvar="gitwatch_enable"
pidfile="/var/run/${name}.pid"
supervisor_pidfile="/var/run/supervisor_${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${gitwatch_enable:="NO"}
: ${gitwatch_user:="nobody"}
: ${gitwatch_group:="nobody"}
: ${gitwatch_args:="${name} --config /usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} -P ${supervisor_pidfile} ${procname} ${gitwatch_args}"

start_precmd=gitwatch_startprecmd

gitwatch_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${gitwatch_user} -g ${gitwatch_group} /dev/null ${pidfile};
    fi

    if [ ! -e ${supervisor_pidfile} ]; then
        install -o ${gitwatch_user} -g ${gitwatch_group} /dev/null ${supervisor_pidfile};
    fi
}

run_rc_command "$1"
