# comms/znet/files/harvest.in
#!/bin/sh

# PROVIDE: harvest
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# harvest_enable (bool):          Set to NO by default.
#               Set it to YES to enable harvest.
# harvest_user (string):          Set user that harvest will run under
#               Default is "nobody".
# harvest_group (string):         Set group that harvest will run under
#               Default is "nobody".
# harvest_args (string):          Set extra arguments to pass to harvest
#               Default is "".

. /etc/rc.subr

name="harvest"
rcvar="harvest_enable"
pidfile="/var/run/${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${harvest_enable:="NO"}
: ${harvest_user:="nobody"}
: ${harvest_group:="nobody"}
: ${harvest_args:="${name} --config /usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} ${procname} ${harvest_args}"

run_rc_command "$1"
