# comms/znet/files/thingharvest.in
#!/bin/sh

# PROVIDE: thingharvest
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# thingharvest_enable (bool):          Set to NO by default.
#               Set it to YES to enable thingharvest.
# thingharvest_user (string):          Set user that thingharvest will run under
#               Default is "nobody".
# thingharvest_group (string):         Set group that thingharvest will run under
#               Default is "nobody".
# thingharvest_args (string):          Set extra arguments to pass to thingharvest
#               Default is "".
# thingharvest_listen_address (string):Set ip:port that thingharvest will listen on
#               Default is ":9100".
# thingharvest_textfile_dir (string):  Set directory that thingharvest will watch
#               Default is "/var/tmp/thingharvest".

. /etc/rc.subr

name=thingharvest
rcvar=thingharvest_enable

load_rc_config $name

: ${thingharvest_enable:="NO"}
: ${thingharvest_user:="nobody"}
: ${thingharvest_group:="nobody"}
: ${thingharvest_args:="thingharvest --config /usr/local/etc/thingharvest.yaml"}
: ${thingharvest_textfile_dir:="/var/tmp/thingharvest"}


pidfile=/var/run/thingharvest.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"
command_args="-S -T ${name} -f -p ${pidfile} /usr/bin/env ${procname} ${thingharvest_args}"

start_precmd=thingharvest_startprecmd

thingharvest_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${thingharvest_user} -g ${thingharvest_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${thingharvest_textfile_dir} ]; then
        install -d -o ${thingharvest_user} -g ${thingharvest_group} -m 1755 ${thingharvest_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
