# comms/znet/files/timer.in
#!/bin/sh

# PROVIDE: timer
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# timer_enable (bool):          Set to NO by default.
#               Set it to YES to enable timer.
# timer_user (string):          Set user that timer will run under
#               Default is "nobody".
# timer_group (string):         Set group that timer will run under
#               Default is "nobody".
# timer_args (string):          Set extra arguments to pass to timer
#               Default is "".

. /etc/rc.subr

name="timer"
rcvar="timer_enable"
pidfile=/var/run/timer.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"
command_args="-S -T ${name} -r -f -p ${pidfile} /usr/bin/env ${procname} ${timer_args}"

load_rc_config $name

: ${timer_enable:="NO"}
: ${timer_user:="nobody"}
: ${timer_group:="nobody"}
: ${timer_args:="timer --config /usr/local/etc/timer.yaml"}

run_rc_command "$1"
