# comms/znet/files/znet.in
#!/bin/sh

# PROVIDE: znet
# REQUIRE: networking
# KEYWORD:
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# znet_enable (bool):          Set to NO by default.
#               Set it to YES to enable znet.
# znet_user (string):          Set user that znet will run under
#               Default is "nobody".
# znet_group (string):         Set group that znet will run under
#               Default is "nobody".
# znet_args (string):          Set extra arguments to pass to znet
#               Default is "".

. /etc/rc.subr

name="znet"
rcvar="znet_enable"
pidfile="/var/run/${name}.pid"
supervisor_pidfile="/var/run/supervisor_${name}.pid"
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/znet"

load_rc_config $name

: ${znet_enable:="NO"}
: ${znet_user:="nobody"}
: ${znet_group:="nobody"}
: ${znet_args:="--config.file=/usr/local/etc/${name}.yaml"}

command_args="-T ${name} -r -f -S -p ${pidfile} -P ${supervisor_pidfile} ${procname} ${znet_args}"

start_precmd=znet_startprecmd

znet_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${znet_user} -g ${znet_group} /dev/null ${pidfile};
    fi

    if [ ! -e ${supervisor_pidfile} ]; then
        install -o ${znet_user} -g ${znet_group} /dev/null ${supervisor_pidfile};
    fi
}

run_rc_command "$1"
