# sysutils/prometheus-nats-exporter/files/nats_exporter.in
#!/bin/sh

# PROVIDE: nats_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# nats_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable nats_exporter.
# nats_exporter_user (string):          Set user that nats_exporter will run under
#               Default is "nobody".
# nats_exporter_group (string):         Set group that nats_exporter will run under
#               Default is "nobody".
# nats_exporter_args (string):          Set extra arguments to pass to nats_exporter
#               Default is "".
# nats_exporter_listen_address (string):Set ip:port that nats_exporter will listen on
#               Default is ":9100".
# nats_exporter_textfile_dir (string):  Set directory that nats_exporter will watch
#               Default is "/var/tmp/nats_exporter".

. /etc/rc.subr

name=prometheus-nats-exporter
rcvar=nats_exporter_enable

load_rc_config $name

: ${nats_exporter_enable:="NO"}
: ${nats_exporter_user:="nobody"}
: ${nats_exporter_group:="nobody"}
: ${nats_exporter_args:=""}
: ${nats_exporter_textfile_dir:="/var/tmp/nats_exporter"}


pidfile=/var/run/nats_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/nats_exporter"
command_args="-p ${pidfile} /usr/bin/env ${procname} \
    ${nats_exporter_args}"

start_precmd=nats_exporter_startprecmd

nats_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${nats_exporter_user} -g ${nats_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${nats_exporter_textfile_dir} ]; then
        install -d -o ${nats_exporter_user} -g ${nats_exporter_group} -m 1755 ${nats_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
