# sysutils/freebsd_exporter/files/freebsd_exporter.in
#!/bin/sh

# PROVIDE: freebsd_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# freebsd_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable freebsd_exporter.
# freebsd_exporter_user (string):          Set user that freebsd_exporter will run under
#               Default is "nobody".
# freebsd_exporter_group (string):         Set group that freebsd_exporter will run under
#               Default is "nobody".
# freebsd_exporter_args (string):          Set extra arguments to pass to freebsd_exporter
#               Default is "".
# freebsd_exporter_listen_address (string):Set ip:port that freebsd_exporter will listen on
#               Default is ":9100".

. /etc/rc.subr

name=freebsd_exporter
rcvar=freebsd_exporter_enable

load_rc_config $name

: ${freebsd_exporter_enable:="NO"}
: ${freebsd_exporter_user:="nobody"}
: ${freebsd_exporter_group:="nobody"}
: ${freebsd_exporter_args:=""}
: ${freebsd_exporter_textfile_dir:="/var/tmp/freebsd_exporter"}


pidfile=/var/run/freebsd_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/freebsd_exporter"
command_args="-S -T ${name} -f -p ${pidfile} /usr/bin/env ${procname} \
    ${freebsd_exporter_args}"

start_precmd=freebsd_exporter_startprecmd

freebsd_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${freebsd_exporter_user} -g ${freebsd_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${freebsd_exporter_textfile_dir} ]; then
        install -d -o ${freebsd_exporter_user} -g ${freebsd_exporter_group} -m 1755 ${freebsd_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
