# sysutils/nodemanager/files/nodemanager.in
#!/bin/sh

# PROVIDE: nodemanager
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# nodemanager_enable (bool):          Set to NO by default.
#               Set it to YES to enable nodemanager.
# nodemanager_user (string):          Set user that nodemanager will run under
#               Default is "root".
# nodemanager_group (string):         Set group that nodemanager will run under
#               Default is "wheel".
# nodemanager_args (string):          Set extra arguments to pass to nodemanager
#               Default is "".
# nodemanager_listen_address (string):Set ip:port that nodemanager will listen on
#               Default is ":9100".

. /etc/rc.subr

name=nodemanager
rcvar=nodemanager_enable

load_rc_config $name

: ${nodemanager_enable:="NO"}
: ${nodemanager_user:="root"}
: ${nodemanager_group:="wheel"}
: ${nodemanager_args:=""}


pidfile=/var/run/nodemanager.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/nodemanager"
command_args="-R 5 -S -T ${name} -f -p ${pidfile} /usr/bin/env PATH=/sbin:/bin:/usr/sbin:/usr/bin:/usr/local/sbin:/usr/local/bin ${procname} \
    ${nodemanager_args}"

start_precmd=nodemanager_startprecmd

nodemanager_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${nodemanager_user} -g ${nodemanager_group} /dev/null ${pidfile};
    fi
}

load_rc_config $name
run_rc_command "$1"
