# sysutils/grafana-agent/files/promtail.in
#!/bin/sh

# PROVIDE: promtail
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# promtail_enable (bool):          Set to NO by default.
#               Set it to YES to enable promtail.
# promtail_user (string):          Set user that promtail will run under
#               Default is "nobody".
# promtail_group (string):         Set group that promtail will run under
#               Default is "nobody".
# promtail_args (string):          Set extra arguments to pass to promtail
#               Default is "".
# promtail_listen_address (string):Set ip:port that promtail will listen on
#               Default is ":9100".
# promtail_textfile_dir (string):  Set directory that promtail will watch
#               Default is "/var/tmp/promtail".

. /etc/rc.subr

name=promtail
rcvar=promtail_enable

load_rc_config $name

: ${promtail_enable:="NO"}
: ${promtail_user:="nobody"}
: ${promtail_group:="nobody"}
: ${promtail_args:=""}
: ${promtail_textfile_dir:="/var/tmp/promtail"}


pidfile=/var/run/promtail.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/promtail"
command_args="-S -T ${name} -p ${pidfile} /usr/bin/env ${procname} \
    ${promtail_args}"

start_precmd=promtail_startprecmd

promtail_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${promtail_user} -g ${promtail_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${promtail_textfile_dir} ]; then
        install -d -o ${promtail_user} -g ${promtail_group} -m 1755 ${promtail_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
