# sysutils/openldap_exporter/files/openldap_exporter.in
#!/bin/sh

# PROVIDE: openldap_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# openldap_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable openldap_exporter.
# openldap_exporter_user (string):          Set user that openldap_exporter will run under
#               Default is "nobody".
# openldap_exporter_group (string):         Set group that openldap_exporter will run under
#               Default is "nobody".
# openldap_exporter_args (string):          Set extra arguments to pass to openldap_exporter
#               Default is "".
# openldap_exporter_listen_address (string):Set ip:port that openldap_exporter will listen on
#               Default is ":9100".
# openldap_exporter_textfile_dir (string):  Set directory that openldap_exporter will watch
#               Default is "/var/tmp/openldap_exporter".

. /etc/rc.subr

name=openldap_exporter
rcvar=openldap_exporter_enable

load_rc_config $name

: ${openldap_exporter_enable:="NO"}
: ${openldap_exporter_user:="nobody"}
: ${openldap_exporter_group:="nobody"}
: ${openldap_exporter_args:=""}
: ${openldap_exporter_listen_address:=":9100"}
: ${openldap_exporter_textfile_dir:="/var/tmp/openldap_exporter"}


pidfile=/var/run/openldap_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/openldap_exporter"
command_args="-p ${pidfile} /usr/bin/env ${procname} \
    --web.listen-address=${openldap_exporter_listen_address} \
    ${openldap_exporter_args}"

start_precmd=openldap_exporter_startprecmd

openldap_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${openldap_exporter_user} -g ${openldap_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${openldap_exporter_textfile_dir} ]; then
        install -d -o ${openldap_exporter_user} -g ${openldap_exporter_group} -m 1755 ${openldap_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
