# astro/openweathermap_exporter/files/openweathermap_exporter.in
#!/bin/sh

# PROVIDE: openweathermap_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# openweathermap_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable openweathermap_exporter.
# openweathermap_exporter_user (string):          Set user that openweathermap_exporter will run under
#               Default is "nobody".
# openweathermap_exporter_group (string):         Set group that openweathermap_exporter will run under
#               Default is "nobody".
# openweathermap_exporter_args (string):          Set extra arguments to pass to openweathermap_exporter
#               Default is "".
# openweathermap_exporter_listen_address (string):Set ip:port that openweathermap_exporter will listen on
#               Default is ":9100".
# openweathermap_exporter_textfile_dir (string):  Set directory that openweathermap_exporter will watch
#               Default is "/var/tmp/openweathermap_exporter".

. /etc/rc.subr

name=openweathermap_exporter
rcvar=openweathermap_exporter_enable

load_rc_config $name

: ${openweathermap_exporter_enable:="NO"}
: ${openweathermap_exporter_user:="nobody"}
: ${openweathermap_exporter_group:="nobody"}
: ${openweathermap_exporter_args:=""}
: ${openweathermap_exporter_textfile_dir:="/var/tmp/openweathermap_exporter"}


pidfile=/var/run/openweathermap_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/openweathermap_exporter"
command_args="-S -T ${name} -f -p ${pidfile} /usr/bin/env ${procname} \
    ${openweathermap_exporter_args}"

start_precmd=openweathermap_exporter_startprecmd

openweathermap_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${openweathermap_exporter_user} -g ${openweathermap_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${openweathermap_exporter_textfile_dir} ]; then
        install -d -o ${openweathermap_exporter_user} -g ${openweathermap_exporter_group} -m 1755 ${openweathermap_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
