# net/goflow_exporter/files/goflow_exporter.in
#!/bin/sh

# PROVIDE: goflow_exporter
# REQUIRE: LOGIN
# KEYWORD: shutdown
#
# Add the following lines to /etc/rc.conf.local or /etc/rc.conf
# to enable this service:
#
# goflow_exporter_enable (bool):          Set to NO by default.
#               Set it to YES to enable goflow_exporter.
# goflow_exporter_user (string):          Set user that goflow_exporter will run under
#               Default is "nobody".
# goflow_exporter_group (string):         Set group that goflow_exporter will run under
#               Default is "nobody".
# goflow_exporter_args (string):          Set extra arguments to pass to goflow_exporter
#               Default is "".
# goflow_exporter_listen_address (string):Set ip:port that goflow_exporter will listen on
#               Default is ":9100".

. /etc/rc.subr

name=goflow_exporter
rcvar=goflow_exporter_enable

load_rc_config $name

: ${goflow_exporter_enable:="NO"}
: ${goflow_exporter_user:="nobody"}
: ${goflow_exporter_group:="nobody"}
: ${goflow_exporter_args:=""}
: ${goflow_exporter_textfile_dir:="/var/tmp/goflow_exporter"}


pidfile=/var/run/goflow_exporter.pid
command="/usr/sbin/daemon"
procname="%%PREFIX%%/bin/goflow_exporter"
command_args="-S -T ${name} -f -p ${pidfile} /usr/bin/env ${procname} \
    ${goflow_exporter_args}"

start_precmd=goflow_exporter_startprecmd

goflow_exporter_startprecmd()
{
    if [ ! -e ${pidfile} ]; then
        install -o ${goflow_exporter_user} -g ${goflow_exporter_group} /dev/null ${pidfile};
    fi
    if [ ! -d ${goflow_exporter_textfile_dir} ]; then
        install -d -o ${goflow_exporter_user} -g ${goflow_exporter_group} -m 1755 ${goflow_exporter_textfile_dir}
    fi
}

load_rc_config $name
run_rc_command "$1"
